# www/www/settings/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
[include]
files=/etc/supervisord/conf.d/*.conf

[program:site]
command         = gunicorn www.wsgi:application -c /Users/zijianwang/future/duodianer/www/www/settings/gunicorn_conf.py
directory       = /Users/zijianwang/future/duodianer/www
user            = zijianwang
autostart       = true
autorestart     = true

environment     = PYTHONPATH="/Users/zijianwang/work/persona/"
stdout_logfile  = /Users/zijianwang/work/persona/logs/gunicorn_supervisor.log
stderr_logfile  = /Users/zijianwang/work/persona/logs/gunicorn_supervisor.log